# systemd/brightness-daemon.service
[Unit]
Description=Display brightness control daemon
After=basic.target

[Service]
ExecStart=brightness_daemon

[Install]
WantedBy=multi-user.target
