# systemd/brightness-daemon.service
[Unit]
Description=Display brightness control daemon

[Service]
ExecStart=brightness-daemon
