# systemd/cpu-powerlimit@.service
[Unit]
Description=Set CPU powerlimit to profile %i

[Service]
ExecStart=cpu_powerlimit --profile=%i
