# systemd/gentoo-sandbox.service
[Unit]
Description=Gentoo sandbox chroot environment

[Service]
Type=notify
ExecStart=gentoo_sandbox --prepare
