# systemd/razer-config@.service
[Unit]
Description=Razer DeathAdder V2 configuration on %I

[Service]
Type=oneshot
ExecStart=razer_config %i
