# systemd/user/clean_bashhistory.service
[Unit]
Description=Cleaning of the Bash history by deduplication

[Service]
ExecStart=clean_bashhistory
