# systemd/user/clean_bashhistory.timer
[Unit]
Description=Cleaning of the Bash history by deduplication

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
